# mqtt-netbase.service
[Unit]
Description=Face Recognition Websocket
After=network-online.target

[Service]
ExecStart=/usr/local/lib/mqttnetalert/netbase.sh
Restart=on-abort

[Install]
WantedBy=multi-user.target
